# application-files/etc/init.d/enable_gpio
#!/bin/sh /etc/rc.common
 
START=90
 
boot() {        
        /usr/bin/lua /etc/boot.lua
}
